# frontend/systemd_frontend.service
[Unit]
Description=LLMxMapReduce Frontend Service  
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/LLMxMapReduce/frontend
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
StandardOutput=append:/home/ubuntu/LLMxMapReduce/frontend/logs/systemd.log
StandardError=append:/home/ubuntu/LLMxMapReduce/frontend/logs/systemd_error.log
SyslogIdentifier=llmx-frontend

[Install]
WantedBy=multi-user.target 